# supervisor/nginx.conf
[program:nginx]
autostart=true
autorestart=true
command=/usr/sbin/nginx
process=(program_name)s
