# supervisor/php7-fpm.conf
[program:php7-fpm]
autostart=true
autorestart=true
command=/usr/sbin/php-fpm7.0
process=(program_name)s
